# systemd/iii-engine.service
[Unit]
Description=iii engine (WebSocket RPC broker + HTTP API)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Environment="PATH=/root/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
WorkingDirectory=/opt/quickstart
ExecStart=/root/.local/bin/iii --config /opt/quickstart/config.yaml
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
